# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:kracked-web]
command=/usr/src/kracked-web/kracked-web
